# holywar-bot.service
[Unit]
Description=Holy War Bot - Automated game bot
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/holy-war-bot
ExecStart=/usr/bin/python3 /home/ubuntu/holy-war-bot/holy_war_bot.py
Restart=always
RestartSec=10
StandardOutput=append:/home/ubuntu/holy-war-bot/bot.log
StandardError=append:/home/ubuntu/holy-war-bot/bot_error.log

# Environment
Environment="DISPLAY=:99"

[Install]
WantedBy=multi-user.target

